# package/debian/postinst
#!/bin/sh
# postinst script for image-postprocesor
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see https://www.debian.org/doc/debian-policy/ or
# the debian-policy package

configure()
{
    systemctl stop wc-readarduino.service
	app_name="arduino-reader"
	mkdir /etc/webcam/$app_name.d -p
	cd /var/lib/webcam/$app_name; 
	pip3 install virtualenv
	virtualenv   --system-site-packages -p python3 venv
	chmod +x /var/lib/webcam/$app_name/venv/bin/activate 
	. /var/lib/webcam/$app_name/venv/bin/activate
	pip3 install -r requirements.txt
	chmod +x /var/lib/webcam/$app_name/bin/wc-readarduino
	chmod +x /var/lib/webcam/$app_name/bin/wc-watchdog
	ln -sf /var/lib/webcam/$app_name/bin/wc-readarduino /usr/local/bin/wc-readarduino
	ln -sf /var/lib/webcam/$app_name/bin/wc-watchdog /usr/local/bin/wc-watchdog
	systemctl daemon-reload
    systemctl enable wc-readarduino.service
}


case "$1" in
    configure)
	configure
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
